# pkg/dangersalad-plasma/dangersalad-plasma.install

post_install () {
    systemctl enable --now bluetooth.service
    # systemctl enable --now NetworkManager.service
    # systemctl enable --now sddm.service
}
