# pkg/dangersalad-plasma/dangersalad-plasma.install
post_upgrade () {
    cp /etc/pulse/default.pa /etc/pulse/default.pa.original
    cp /etc/pulse/default.pa.ds /etc/pulse/default.pa
}

post_install () {
    systemctl enable --now bluetooth.service
    post_upgrade
}
